# core/com.claude-stack.strategist.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude-stack.strategist</string>
    <key>ProgramArguments</key>
    <array>
        <string>/home/user/claude-stack/.venv/bin/python</string>
        <string>-u</string>
        <string>-m</string>
        <string>core.revenue_strategist_loop</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/home/user/claude-stack</string>
    <key>StandardOutPath</key>
    <string>/home/user/Library/Logs/claude-stack/strategist.log</string>
    <key>StandardErrorPath</key>
    <string>/home/user/Library/Logs/claude-stack/strategist.err.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/home/user</string>
        <key>PYTHONPATH</key>
        <string>/home/user/claude-stack</string>
        <key>CLAUDE_STACK_DIR</key>
        <string>/home/user/claude-stack</string>
        <key>STRATEGIST_ENABLED</key>
        <string>1</string>
    </dict>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>20</integer>
    </dict>
    <key>RunAtLoad</key>
    <false/>
    <key>AbandonProcessGroup</key>
    <true/>
</dict>
</plist>
